# post-install/nvidia-fan-control.timer
[Unit]
Description=Run nvidia-fan-control every 5 sec

[Timer]
OnBootSec=10min
OnUnitActiveSec=5sec

[Install]
WantedBy=timers.target
